# DistilBERT/network-detector.service
[Unit]
Description=Network Attack Detector
After=network.target

[Service]
Type=simple
User=jvf
WorkingDirectory=/home/jvf/testes/DistilBERT
Environment=PATH=/home/jvf/testes/DistilBERT/venv/bin
ExecStart=/home/jvf/testes/DistilBERT/venv/bin/python /home/jvf/testes/DistilBERT/realtime_network_monitor.py --interactive
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
